# dot_config/systemd/user/hyprsunset.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/hyprsunset -t 6500
