# dot_config/systemd/user/paccache-local.timer
[Unit]
Description=Discard local AUR package archives weekly

[Timer]
OnCalendar=weekly
AccuracySec=1h
Persistent=true

[Install]
WantedBy=timers.target
