# dot_config/systemd/user/paccache-local.service
[Unit]
Description=Remove cached local AUR package archives.

[Service]
Type=oneshot
ExecStart=/usr/bin/paccache -c ${HOME}/.cache/pacman/pkg -rk0
