# dot_config/systemd/user/private_cliphist-wipe.service
[Unit]
Description=Clear clipboard

[Service]
Type=oneshot
ExecStart=/usr/bin/cliphist wipe

