# dot_config/systemd/user/private_purge-src.service
[Unit]
Description=Clear ~/.local/src directory 

[Service]
Type=oneshot
ExecStart=/usr/bin/sh -c 'rm -rf "${HOME}/.local/src/"*'
